# etc/profile.d/python36.sh
#!/bin/bash

source /opt/rh/rh-python36/enable
export X_SCLS="`scl enable rh-python36 'echo $X_SCLS'`"
